# gaggia-bot.service
[Unit]
Description=GaggiaMate Shot Intelligence Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=REPLACE_WITH_YOUR_USER
WorkingDirectory=REPLACE_WITH_INSTALL_PATH
EnvironmentFile=REPLACE_WITH_INSTALL_PATH/.env
ExecStart=REPLACE_WITH_INSTALL_PATH/.venv/bin/python main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gaggia-bot

[Install]
WantedBy=multi-user.target
