# gaggia-bot.service
[Unit]
Description=GaggiaMate Shot Intelligence Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=michael
WorkingDirectory=/home/michael/gaggia-bot
EnvironmentFile=/home/michael/gaggia-bot/.env
ExecStart=/home/michael/gaggia-bot/.venv/bin/python main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
